# thbigmatome/scripts/backup/cron.d/backup
# thbigmatome 本番DB自動バックアップ cron設定
#
# デプロイ方法（VPS上でdeployユーザーとして）:
#   crontab -e  で以下の行を追加
#   または: crontab scripts/backup/cron.d/backup
#
# ⚠️ このファイルは crontab 形式（/etc/cron.d 形式ではない）
#    ユーザー名フィールドは不要

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# 毎日 02:00 (JST) にバックアップ実行
0 2 * * * /home/deploy/projects/thbigmatome/scripts/backup/backup_daily.sh >> /home/deploy/backups/thbigmatome/cron.log 2>&1
